# iwango.service.in
[Unit]
Description=IWANGO gate and lobby server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=INSTALL_USER
ExecStart=SBINDIR/iwango_server SYSCONFDIR/iwango.cfg
StandardOutput=append:/var/log/iwango.log

[Install]
WantedBy=multi-user.target
